# kafka/etc/supervisor/conf.d/kafka.conf
[program:kafka]
command=/opt/kafka/kafka-start.sh
user=fitbur
directory=/tmp
autostart=true
autorestart=true
stdout_logfile=/out/kafka_out.log
stderr_logfile=/out/kafka_err.log
