# images/base/datastax/datastax-opscenter/etc/supervisor/conf.d/opscenter.conf
[program:opscenter]
command=/opt/opscenter/opscenter-start.sh
user=fitbur
directory=/tmp
autostart=true
autorestart=true
stdout_logfile=/out/opscenter_out.log
stderr_logfile=/out/opscenter_err.log
