# docker/storm/supervisor/etc/supervisor/conf.d/storm-supervisor.conf
[program:storm-supervisor]
command=/usr/bin/storm supervisor
directory=/out
autorestart=true
user=nobody
stdout_logfile=/out/supervisor.out.log
stderr_logfile=/out/supervisor.err.log

