# docker/storm/supervisor/etc/supervisor/conf.d/storm-logviewer.conf
[program:storm-logviewer]
command=/usr/bin/storm logviewer
directory=/out
autorestart=true
user=nobody
stdout_logfile=/out/logviewer.out.log
stderr_logfile=/out/logviewer.err.log
