# docker/storm/nimbus/etc/supervisor/conf.d/storm-nimbus.conf
[program:storm-nimbus]
command=/usr/bin/storm nimbus
directory=/home/fitbur
autorestart=true
user=fitbur
stdout_logfile=/out/nimbus.out.log
stderr_logfile=/out/nimbus.err.log
