# docker/storm/nimbus/etc/supervisor/conf.d/storm-drpc.conf
[program:storm-drcp]
command=/usr/bin/storm drpc
directory=/home/fitbur
autorestart=true
user=fitbur
stdout_logfile=/out/drpc.out.log
stderr_logfile=/out/drpc.err.log
