# nexus/etc/supervisor/conf.d/nexus.conf
[program:nexus]
command=/opt/nexus/nexus-start.sh
user=fitbur
directory=/tmp
autostart=true
autorestart=true
stdout_logfile=/out/nexus_out.log
stderr_logfile=/out/nexus_err.log
