# cassandra/etc/supervisor/conf.d/cassandra.conf
[program:cassandra]
command=/opt/cassandra/cassandra-start.sh
user=fitbur
directory=/tmp
autostart=true
autorestart=true
stdout_logfile=/out/cassandra_out.log
stderr_logfile=/out/cassandra_err.log
