# images/zookeeper/etc/supervisor/conf.d/zookeeper.conf
[program:zookeeper]
command=/opt/zookeeper/zookeeper-start.sh
user=fitbur
directory=/out
autostart=true
autorestart=true
stdout_logfile=/out/zookeeper_out.log
stderr_logfile=/out/zookeeper_err.log
