# images/zookeeper/etc/supervisor/conf.d/zookeeper.conf
[program:zookeeper]
command=/opt/zookeeper/zookeeper-start.sh
directory=/tmp
autostart=true
autorestart=true
stdout_logfile=/out/supervisor_out.log
stderr_logfile=/out/supervisor_err.log
