# base/common/ubuntu/etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
directory=/out
autostart=true
autorestart=true
