# docker/storm/ui/etc/supervisor/conf.d/storm-ui.conf
[program:storm-ui]
command=/usr/bin/storm ui
directory=/home/fitbur
autorestart=true
user=fitbur
stdout_logfile=/out/ui.out.log
stderr_logfile=/out/ui.err.log
