# base/datastax/datastax-agent/etc/supervisor/conf.d/agent.conf
[program:agent]
command=/opt/datastax-agent/agent-start.sh
user=fitbur
directory=/tmp
autostart=true
autorestart=true
stdout_logfile=/out/agent_out.log
stderr_logfile=/out/agent_err.log
