# jenkins/etc/supervisor/conf.d/jenkins.conf
[program:jenkins]
command=/opt/jenkins/jenkins-start.sh
user=fitbur
directory=/tmp
autostart=true
autorestart=true
stdout_logfile=/out/jenkins_out.log
stderr_logfile=/out/jenkins_err.log
